# devops/pollers/home-poller/supervisord.conf
; StakTrakr Home Poller — Supervisor Configuration
; Manages cron, dashboard, and metrics exporter.
; Unlike the Fly.io container, no embedded Firecrawl/Redis/RabbitMQ/PostgreSQL.

; ── Cron (retail poller, spot poller, provider export, goldback) ──

[program:cron]
command=cron -f
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Dashboard (provider editor + status UI) ──

[program:dashboard]
command=node /app/dashboard.js
directory=/app
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Prometheus Metrics Exporter ──

[program:metrics-exporter]
command=node /app/metrics-exporter.js
directory=/app
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
